# config/systemd/library-heartbeat.service
[Unit]
Description=Library Heartbeat — continuous lane health beacon
After=network.target

[Service]
Type=simple
ExecStart=/home/we4free/.nvm/versions/node/v20.20.2/bin/node /home/we4free/agent/repos/self-organizing-library/scripts/heartbeat.js --lane library --continuous --interval 60
WorkingDirectory=/home/we4free/agent/repos/self-organizing-library
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=LANE_NAME=library
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus

[Install]
WantedBy=default.target
